# systemd/infra-maintenance.timer
[Unit]
Description=Run Maintenance Daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
